# spook.service
[Unit]
Description=Deer spooking as a service

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/spook.sh

[Install]
WantedBy=multi-user.target
